# docker/inotifyd-init
#!/sbin/openrc-run

command=/sbin/inotifyd
command_args="$INOTIFYD_ARGS"
command_user="$INOTIFYD_USER"
pidfile=/run/${RC_SVCNAME}.pid
command_background=yes
start_stop_daemon_args="--stdout /var/log/$RC_SVCNAME/${RC_SVCNAME}.log --stderr /var/log/$RC_SVCNAME/${RC_SVCNAME}.log"

start_pre() {
  checkpath --directory --owner $INOTIFYD_USER --mode 0775 /var/log/$RC_SVCNAME
}